# infra/quadlet/tavianator.com-blog.container
[Unit]
Description=tavianator.com-blog

[Container]
Image=docker.io/tavianator/tavianator.com-blog
AutoUpdate=registry
Pod=tavianator.com.pod
