# certbot.timer
[Unit]
Description=Let's Encrypt autorenewal

[Timer]
OnCalendar=0/12:00:00
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target
